# services/systemd/xwd.service.in
[Unit]
Description=XWallpaper daemon

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/xwd
Restart=on-failure

[Install]
WantedBy=default.target
